# roles/game_server/templates/zerotier.service.j2
[Unit]
Description=ZeroTier One Service Daemon

[Service]
Type=forking
User=root
KillMode=none
ExecStart={{ gs_zt_which['stdout_lines'][0] }} -d

[Install]
WantedBy=multi-user.target